# systemd/openclaw-skillbank.timer
[Unit]
Description=Periodic openclaw-skillbank distillation + export

[Timer]
# Once nightly. Distillation isn't time-sensitive — running it during quiet
# hours keeps it well clear of any user-facing agent traffic and trivially
# avoids rate-limit contention even on shared buckets.
# Override by dropping a unit file at ~/.config/systemd/user/openclaw-skillbank.timer.d/override.conf
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=30min
Persistent=true

[Install]
WantedBy=timers.target
